# deploy/ansible/roles/agent/templates/wraith-agent.service.j2
[Unit]
Description=Wraith Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ wraith_user }}
Group={{ wraith_group }}
WorkingDirectory={{ wraith_install_dir }}
EnvironmentFile=-{{ wraith_agent_env_file }}
ExecStart={{ wraith_install_dir }}/bin/agent --config {{ wraith_config_dir }}/agent.json
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ wraith_log_dir }} {{ wraith_config_dir }}

[Install]
WantedBy=multi-user.target
