# deploy/ansible/roles/control_plane/templates/wraith-server.service.j2
[Unit]
Description=Wraith Control Plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ wraith_user }}
Group={{ wraith_group }}
WorkingDirectory={{ wraith_install_dir }}
EnvironmentFile={{ wraith_server_env_file }}
ExecStart={{ wraith_install_dir }}/bin/server --grpc-addr {{ wraith_server_grpc_addr }} --enrollment-addr {{ wraith_server_enrollment_addr }} --mtls-endpoint {{ wraith_server_mtls_endpoint }} --catalog {{ wraith_config_dir }}/agents.json --enrollment-tokens {{ wraith_config_dir }}/enrollment-tokens.json --ca-cert {{ wraith_cert_dir }}/ca.crt --ca-key {{ wraith_cert_dir }}/ca.key --server-cert {{ wraith_cert_dir }}/server.crt --server-key {{ wraith_cert_dir }}/server.key --database-url {{ wraith_database_url }}{% if wraith_server_offline_ca %} --offline-ca{% endif %}
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ wraith_data_dir }} {{ wraith_log_dir }} {{ wraith_config_dir }}

[Install]
WantedBy=multi-user.target
